# container-files/etc/supervisor.d/smbd.conf
[program:smbd]
command=/usr/sbin/smbd -D
autostart=true
stderr_logfile = NONE
stdout_logfile = /var/log/secure
redirect_stderr = true
