# container-files/etc/cron.d/storage-sync
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=kanalfred@gmail.com

### Rsnyc ###
# rsync backup storage - Daily every day after 5:00am
00 05 * * * root /config/scripts/storage-sync.sh
